# core/services.d/websockify/run
#!/usr/bin/with-contenv bash

/opt/websockify/run $TARGET_PORT localhost:5900