# core/services.d/xvfb/run
#!/usr/bin/with-contenv bash

rm -rf /tmp/.X${DISPLAY}-lock
/usr/bin/Xvfb ${DISPLAY} -screen ${SCREEN} ${SCREEN_WIDTH}x${SCREEN_HEIGHT}x${SCREEN_DEPTH}