# core/services.d/backend/run
#!/usr/bin/with-contenv bash

node ./build/app.js